# templates/testodoo.timer
[Unit]
Description=Reiniciar docker de test

[Timer]
OnBootSec=15min
OnCalendar=*-*-* 08:00:00
Persistent=true

[Install]
WantedBy=timers.target
